# root/etc/services.d/jenkins/run
#!/usr/bin/with-contenv bash
exec s6-setuidgid cont java -jar /app/jenkins.war $JENKINS_OPTS